# external/nginx/certbot.cron
# Let's Encrypt renewal
05 5 * * * certbot renew --renew-hook "systemctl reload nginx"